# target/usr/lib/systemd/system/battery.target
[Unit]
Description=On battery
DefaultDependencies=no
StopWhenUnneeded=yes
